# files/uuplugin.init
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
	rm -rf /tmp/uu && mkdir /tmp/uu && cd /tmp/uu
	procd_open_instance
	procd_set_param command /usr/sbin/uu/uuplugin /etc/uu.conf
	procd_close_instance
}
